# etc/cron.d/phoenix-stretch
* * * * * root /usr/bin/reprepro --silent --silent -b /data/reprepro/xivo --waitforlock 10 --noskipold update phoenix-stretch
