# etc/cron.d/phoenix-stretch
* * * * * root /usr/local/bin/reprepro-main --silent --silent --waitforlock 10 --noskipold update phoenix-stretch
